# claw-control.service
[Unit]
Description=Claw Control Dashboard (Next.js)
After=network.target claw-memory.service

[Service]
Type=simple
User=mat
WorkingDirectory=/home/mat/Projects/claw-control
ExecStartPre=/bin/sh -c "fuser -k 3030/tcp || true"
ExecStart=/home/linuxbrew/.linuxbrew/bin/npx next start -H 0.0.0.0 -p 3030
Restart=always
RestartSec=5
KillMode=control-group
TimeoutStopSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
