# claw-control.service
[Unit]
Description=Claw Control Dashboard (Next.js)
After=network.target claw-memory.service

[Service]
Type=simple
User=mat
WorkingDirectory=/home/mat/Projects/claw-control
ExecStart=/home/linuxbrew/.linuxbrew/bin/npx next start -H 0.0.0.0 -p 3030
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
